# deploy/systemd/csd-worker.service
[Unit]
Description=Copilot Studio Deployer - Background Worker
Documentation=https://github.com/your-org/copilot-studio-deployer
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=csd
Group=csd
WorkingDirectory=/opt/copilot-studio-deployer/packages/worker
ExecStart=/usr/bin/node /opt/copilot-studio-deployer/packages/worker/dist/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=csd-worker

# Environment
Environment=NODE_ENV=production
EnvironmentFile=/etc/csd/environment

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/copilot-studio-deployer/snapshots
ReadWritePaths=/opt/copilot-studio-deployer/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
